# tarsierpos-backend.service
[Unit]
Description=TarsierPos Backend (Django :9000)
After=network.target

[Service]
User=__DEPLOY_USER__
WorkingDirectory=/home/__DEPLOY_USER__/TarsierPOS
ExecStart=/home/__DEPLOY_USER__/TarsierPOS/venv/bin/gunicorn pos_config.wsgi:application --bind 127.0.0.1:9000 --workers 3 --timeout 60
ExecStopPost=/bin/bash /home/__DEPLOY_USER__/TarsierPOS/backup_db.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
