# tarsierpos-frontend.service
# DEPRECATED — frontend is now served by nginx over HTTPS (FIX-047)
# nginx site config: /etc/nginx/sites-available/canteen-pos
# TLS cert: /home/ralph/TarsierPOS/certs/ralph-optiplex-3060.tail037b5f.ts.net.{crt,key}
# Cert renewal: sudo tailscale cert ralph-optiplex-3060.tail037b5f.ts.net
#               then: sudo mv ralph-optiplex-3060.*.{crt,key} /home/ralph/TarsierPOS/certs/
#               then: sudo systemctl reload nginx
#
# This Python HTTP server service is DISABLED.
# To re-enable plain-HTTP fallback on :8081 (local/dev only), restore this file.
[Unit]
Description=TarsierPos Frontend — DISABLED (nginx handles HTTPS, see FIX-047)
After=network.target
ConditionPathExists=/dev/null

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
