# scripts/systemd/tarsierpos.service.d/ordering.conf
# /etc/systemd/system/tarsierpos.service.d/ordering.conf
# FEATURE-022: wait for tailscale + a usable network before starting
# Install: sudo install -m 0644 -o root -g root \
#   scripts/systemd/tarsierpos.service.d/ordering.conf \
#   /etc/systemd/system/tarsierpos.service.d/ordering.conf
# Then: sudo systemctl daemon-reload
[Unit]
After=tailscaled.service network-online.target
Wants=network-online.target
