# scripts/systemd/tarsierpos.service.d/preload.conf
# /etc/systemd/system/tarsierpos.service.d/preload.conf
# FLAG-067: --preload catches import errors at service start, not first request
# Install: sudo mkdir -p /etc/systemd/system/tarsierpos.service.d/
#          sudo install -m 0644 -o root -g root scripts/systemd/tarsierpos.service.d/preload.conf /etc/systemd/system/tarsierpos.service.d/preload.conf
#          sudo systemctl daemon-reload
[Service]
ExecStart=
ExecStart=/home/ralph/TarsierPOS/venv/bin/gunicorn pos_config.wsgi:application --bind 127.0.0.1:9000 --workers 3 --timeout 60 --preload
