# config/systemd/tarsierpos.service.d/limits.conf
[Unit]
StartLimitIntervalSec=600
StartLimitBurst=10

[Service]
StartLimitBurst=10
