# scripts/systemd/tarsierpos-backup-offsite.service.d/gcs.conf
# /etc/systemd/system/tarsierpos-backup-offsite.service.d/gcs.conf
# FEATURE-017(b): GCS upload via service account
# Install: sudo install -m 0644 -o root -g root \
#   scripts/systemd/tarsierpos-backup-offsite.service.d/gcs.conf \
#   /etc/systemd/system/tarsierpos-backup-offsite.service.d/gcs.conf
# Then: sudo systemctl daemon-reload
[Service]
Environment=GOOGLE_APPLICATION_CREDENTIALS=/etc/tarsierpos/gcs-backup-sa.json
Environment=GCS_DEST=gs://tarsierpos-backups-rma/dev-optiplex
